# supervisor/axiom-api.conf
[program:axiom-api]
command=/opt/axiom/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 2 --proxy-headers
directory=/opt/axiom/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/axiom/api.log
stderr_logfile=/var/log/axiom/api.err.log
environment=PYTHONUNBUFFERED="1"
user=axiom
