# supervisor/axiom-beat.conf
[program:axiom-beat]
command=/opt/axiom/venv/bin/celery -A app.celery_app beat --loglevel=info
directory=/opt/axiom/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/axiom/beat.log
stderr_logfile=/var/log/axiom/beat.err.log
environment=PYTHONUNBUFFERED="1"
user=axiom
