# supervisor/axiom-worker.conf
[program:axiom-worker]
command=/opt/axiom/venv/bin/celery -A app.celery_app worker --loglevel=info --concurrency=4
directory=/opt/axiom/app
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stopwaitsecs=30
stdout_logfile=/var/log/axiom/worker.log
stderr_logfile=/var/log/axiom/worker.err.log
environment=PYTHONUNBUFFERED="1"
user=axiom
