# supervisor/axiom-feed.conf
[program:axiom-feed]
command=/opt/axiom/venv/bin/python -m app.feed
directory=/opt/axiom/app
autostart=true
autorestart=true
startretries=20
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/axiom/feed.log
stderr_logfile=/var/log/axiom/feed.err.log
environment=PYTHONUNBUFFERED="1"
user=axiom
